# supervisord.conf
# Supervisor config — manages the long-running app processes.
# Postgres / Redis / MinIO are started (and migrated) by the entrypoint BEFORE
# supervisord launches, so the app processes always find a ready backend.

[supervisord]
nodaemon=true
logfile=/data/app/supervisord.log
pidfile=/data/app/supervisord.pid
user=root
# Preserve the environment exported by the entrypoint (REDIS_URL, DATABASE_URL,
# STORAGE_DRIVER, etc.) so child processes inherit all runtime configuration.
clear_environment=false

[program:web]
command=node server.js
directory=/app/apps/web
priority=30
autostart=true
autorestart=true
startsecs=8
stopwaitsecs=20
stdout_logfile=/data/app/web.out.log
stderr_logfile=/data/app/web.err.log
environment=PORT="7860",HOSTNAME="0.0.0.0",REDIS_URL="%(ENV_REDIS_URL)s",DATABASE_URL="%(ENV_DATABASE_URL)s",STORAGE_DRIVER="%(ENV_STORAGE_DRIVER)s",STORAGE_LOCAL_DIR="%(ENV_STORAGE_LOCAL_DIR)s",OCR_PROVIDER="%(ENV_OCR_PROVIDER)s",OCR_CLOUD_ENABLED="%(ENV_OCR_CLOUD_ENABLED)s",QUEUE_DRIVER="%(ENV_QUEUE_DRIVER)s",NODE_ENV="production"

[program:ocr-worker]
command=node --import tsx workers/ocr-worker/src/index.ts
directory=/app
priority=31
autostart=true
autorestart=true
startsecs=8
stopwaitsecs=20
stdout_logfile=/data/app/ocr-worker.out.log
stderr_logfile=/data/app/ocr-worker.err.log
environment=REDIS_URL="%(ENV_REDIS_URL)s",DATABASE_URL="%(ENV_DATABASE_URL)s",STORAGE_DRIVER="%(ENV_STORAGE_DRIVER)s",STORAGE_LOCAL_DIR="%(ENV_STORAGE_LOCAL_DIR)s",OCR_PROVIDER="%(ENV_OCR_PROVIDER)s",OCR_CLOUD_ENABLED="%(ENV_OCR_CLOUD_ENABLED)s",QUEUE_DRIVER="%(ENV_QUEUE_DRIVER)s",NODE_ENV="production"

[program:export-worker]
command=node --import tsx workers/export-worker/src/index.ts
directory=/app
priority=32
autostart=true
autorestart=true
startsecs=8
stopwaitsecs=20
stdout_logfile=/data/app/export-worker.out.log
stderr_logfile=/data/app/export-worker.err.log
environment=REDIS_URL="%(ENV_REDIS_URL)s",DATABASE_URL="%(ENV_DATABASE_URL)s",STORAGE_DRIVER="%(ENV_STORAGE_DRIVER)s",STORAGE_LOCAL_DIR="%(ENV_STORAGE_LOCAL_DIR)s",OCR_PROVIDER="%(ENV_OCR_PROVIDER)s",OCR_CLOUD_ENABLED="%(ENV_OCR_CLOUD_ENABLED)s",QUEUE_DRIVER="%(ENV_QUEUE_DRIVER)s",NODE_ENV="production"
