# supervisord.conf
# Supervisor config — manages the long-running app processes.
# Postgres / Redis / MinIO are started (and migrated) by the entrypoint BEFORE
# supervisord launches, so the app processes always find a ready backend.

[supervisord]
nodaemon=true
logfile=/data/app/supervisord.log
pidfile=/data/app/supervisord.pid
user=root

[program:web]
command=node apps/web/server.js
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=8
stopwaitsecs=20
stdout_logfile=/data/app/web.out.log
stderr_logfile=/data/app/web.err.log
environment=PORT="7860",HOSTNAME="0.0.0.0"

[program:ocr-worker]
command=node --import tsx workers/ocr-worker/src/index.ts
directory=/app
priority=31
autostart=true
autorestart=true
startsecs=8
stopwaitsecs=20
stdout_logfile=/data/app/ocr-worker.out.log
stderr_logfile=/data/app/ocr-worker.err.log

[program:export-worker]
command=node --import tsx workers/export-worker/src/index.ts
directory=/app
priority=32
autostart=true
autorestart=true
startsecs=8
stopwaitsecs=20
stdout_logfile=/data/app/export-worker.out.log
stderr_logfile=/data/app/export-worker.err.log
